# tmux/persist/tmux-persist.cron
# tmux-persistence

* * * * *  dosborne  cd /home/dosborne/.dotfiles/tmux/persist && ./tmux-persist.rb  
